# scripts/afd.in
#!/bin/sh
#
# afd:       Starts the AFD
#
# chkconfig: 235 70 10
# description: AFD is an automatic file distributor and collector. It is \
#              used to distribute data files via FTP, FTPS, SFTP, SMTP \
#              and/or SCP. It collects files via FTP, FTPS, SFTP, HTTP \
#              and/or HTTPS.
### BEGIN INIT INFO
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# Default-Start: 2 3 5
# Default-Stop: 0 1 6
# Short-Description: Starts the AFD
# Description: AFD is an automatic file distributor and collector. It is
#              used to distribute data files via FTP, FTPS, SFTP, SMTP
#              and/or SCP. It collects files via FTP, FTPS, SFTP, HTTP
#              and/or HTTPS.
### END INIT INFO


CURRENT_SCRIPT_NAME="`basename $0`"
if [ -r @SYSCONFIGDIR@/afd ]
then
   if [ "$CURRENT_SCRIPT_NAME" = "afdha" ]
   then
      AFD_USERS="`grep \^AFD_HA_USER_ @SYSCONFIGDIR@/afd | cut -d '=' -f2`"
   else
      AFD_USERS="`grep \^AFD_USER_ @SYSCONFIGDIR@/afd | cut -d '=' -f2`"
   fi
else
   AFD_USERS=""
fi


# See how we were called.
case "$1" in
  start)
	RETVAL=1
	for afduser in $AFD_USERS
	do
	   echo -n "Starting AFD for $afduser : "
	   su - $afduser -c "afd -a"
	   if [ $? -eq 0 ]
	   then
	      RETVAL=0
	      touch /var/lock/subsys/afd.$afduser
	      echo "Done"
	   else
	      echo "Failed"
	   fi
	done
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/afd
	;;
  stop)
	RETVAL=1
	for afduser in $AFD_USERS
	do
	   echo -n "Stopping AFD for $afduser : "
	   su - $afduser -c "afd -S"
	   case $? in
	      0)
	         RETVAL=0
	         rm -f /var/lock/subsys/afd.$afduser
	         echo "Done"
	         ;;
	      10)
	         RETVAL=0
	         rm -f /var/lock/subsys/afd.$afduser
	         echo "Is already stopped"
	         ;;
	      *)
	         echo "Failed"
	         ;;
	   esac
	done
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/afd
	;;
  restart)
	$0 stop
	sleep 2
	$0 start
	RETVAL=$?
	;;
  condrestart)
	RETVAL=1
	for afduser in $AFD_USERS
	do
	   if [ -f /var/lock/subsys/restart_afd.$afduser ]
	   then
	      echo -n "Starting AFD for $afduser : "
	      su - $afduser -c "afd -a"
	      if [ $? -eq 0 ]
	      then
	         RETVAL=0
	         touch /var/lock/subsys/afd.$afduser
	         echo "Done"
	      else
	         echo "Failed"
	      fi
	      rm -f /var/lock/subsys/restart_afd.$afduser
	   else
	      su - $afduser -c "afd -S"
	      if [ $? -eq 0 ]
	      then
	         sleep 2
	         su - $afduser -c "afd -a"
	         if [ $? -eq 0 ]
	         then
	            RETVAL=0
	            touch /var/lock/subsys/afd.$afduser
	            echo "Done"
	         else
	            echo "Failed"
	         fi
	      fi
	   fi
	done

	if [ ! "$CURRENT_SCRIPT_NAME" = "afdha" ]
	then
	   # HA users may only be started if we find a restart_afd.user file!
	   if [ -r @SYSCONFIGDIR@/afd ]
	   then
	      AFD_HA_USERS="`grep \^AFD_HA_USER_ @SYSCONFIGDIR@/afd | cut -d '=' -f2`"
	   fi
	   for afduser in $AFD_HA_USERS
	   do
	      if [ -f /var/lock/subsys/restart_afd.$afduser ]
	      then
	         echo -n "Starting AFD for $afduser : "
	         su - $afduser -c "afd -a"
	         if [ $? -eq 0 ]
	         then
	            RETVAL=0
	            touch /var/lock/subsys/afd.$afduser
	            echo "Done"
	         else
	            echo "Failed"
	         fi
	         rm -f /var/lock/subsys/restart_afd.$afduser
	      fi
	   done
	fi
	rm -f /var/lock/subsys/afd.restart
	;;
  status)
	NOT_RUNNING=""
	IS_RUNNING=""
	for afduser in $AFD_USERS
	do
	   su - $afduser -c "afd -c"
	   RETVAL=$?
	   if [ $RETVAL -ne 5 ]
	   then
	      NOT_RUNNING="`echo $NOT_RUNNING` $afduser"
	   else
	      IS_RUNNING="`echo $IS_RUNNING` $afduser"
	   fi
	done
	if [ "$NOT_RUNNING" = "" ]
	then
	   echo "AFD for $IS_RUNNING is running."
	   RETVAL=0
	else
	   echo "AFD for $NOT_RUNNING is NOT active."
	   RETVAL=1
	fi
	;;
  *)
	echo "*** Usage: afd {start|stop|restart|condrestart|status}"
	RETVAL=1
	;;
esac

exit $RETVAL
